# auto-theme/theme-reschedule.timer
[Unit]
Description=Auto theme switcher — reschedule at midnight

[Timer]
OnCalendar=*-*-* 00:01:00
Persistent=true

[Install]
WantedBy=timers.target
